# src/TechStacks/scripts/server/supervisord/web-prerender.conf
[program:web-prerender]
command=node /home/deploy/apps/scripts/prerender.js
autostart=true
autorestart=true
environment=NODE_ENV=production
stderr_logfile=/var/log/prerender.err.log
stdout_logfile=/var/log/prerender.out.log
user=ubuntu
stopsignal=INT
