# src/TechStacks/scripts/server/supervisord/web-techstacks.conf
[program:web-techstacks]
command=/usr/bin/dotnet /home/deploy/apps/techstacks/TechStacks.dll
directory=/home/deploy/apps/techstacks
autostart=true
autorestart=true
stderr_logfile=/var/log/web-techstacks.err.log
stdout_logfile=/var/log/web-techstacks.out.log
environment=ASPNETCORE_ENVIRONMENT=Production,ASPNETCORE_URLS="http://*:16325/",SERVICESTACK_LICENSE="",FORUMS_DB=""
user=deploy
stopsignal=INT